# tmp/amica.service
[Unit]
Description=Amica Web Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=purestory
Group=purestory
WorkingDirectory=/home/purestory/amica
Environment=NODE_ENV=production
Environment=PORT=3100
Environment=BASE_PATH=/amica
ExecStartPre=/usr/local/bin/fix-manifest.sh
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=amica

[Install]
WantedBy=multi-user.target 